# scripts/sync_script.cron
* * * * * bash -c '/scripts/sync_script.sh'
